# debian/vyatta-wanloadbalance.postrm.in
#!/bin/bash
which /etc/init.d/vyatta-wanloadbalance >/dev/null && /etc/init.d/vyatta-wanloadbalance stop

update-rc.d vyatta-wanloadbalance remove

