# debian/vyatta-wanloadbalance.postinst.in
#!/bin/bash

LOADBALANCE_DIR=/var/load-balance


prefix=/opt/vyatta
sysconfdir=${prefix}/etc

mkdir -p $LOADBALANCE_DIR
chmod 0755 $LOADBALANCE_DIR

update-rc.d -f vyatta-wanloadbalance defaults 2>/dev/null

#now move op cmds to correct location
mv /opt/vyatta/share/vyatta-cfg/templates/show/wan-load-balance /opt/vyatta/share/vyatta-op/templates/show/wan-load-balance